# HOME/.config/autostart/cmst.desktop
[Desktop Entry]
Exec=cmst -m -w 5
Name=cmst
Type=Application
Version=1.0
