# HOME/.config/autostart/picom.desktop
[Desktop Entry]
Exec=/home/takase/.config/manage_picom start
Name=picom
Type=Application
Version=1.0
