# HOME/.config/autostart/wallpaper.desktop
[Desktop Entry]
Exec=/home/takase/.config/manage_wp start
Name=Wallpaper
Type=Application
Version=1.0
