# HOME/.config/autostart/xautolock.desktop
[Desktop Entry]
Exec=xautolock
Name=xautolock
Type=Application
Version=1.0
