# dot_config/autostart/executable_restore-wallpaper.desktop
[Desktop Entry]
Type=Application
Exec=nitrogen --restore
Terminal=false
Hidden=false
Name=Restore wallpaper
