# dot_config/autostart/executable_dunst.desktop
[Desktop Entry]
Type=Application
Exec=dunst
Terminal=false
Hidden=false
Name=dunst
