# dot_config/autostart/executable_picom.desktop
[Desktop Entry]
Type=Application
Exec=picom
Terminal=false
Hidden=false
Name=picom
