# dot_config/autostart/executable_plank.desktop
[Desktop Entry]
Type=Application
Exec=plank
Terminal=false
Hidden=false
Name=plank
