# dot_config/autostart/executable_xfce4-power-manager.desktop
[Desktop Entry]
Type=Application
Exec=xfce4-power-manager
Terminal=false
Hidden=false
Name=xfce4-power-manager
