# dot_config/autostart/executable_polybar.desktop
[Desktop Entry]
Type=Application
Exec=bash ~/.config/openbox/polybar/launch.sh
Terminal=false
Hidden=false
Name=polybar
